# todo-recv.service
[Unit]
Description=飞书电话提醒系统 - 消息接收器
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/todo-voice-caller
ExecStart=/usr/bin/python3.12 /home/ubuntu/todo-voice-caller/recv.py
Restart=always
RestartSec=5
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="HOME=/home/ubuntu"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
